# m3u8-emitter.service
[Unit]
Description=M3U8 to RTMP Emitter Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/m3u8-emisor
Environment=NODE_ENV=production
Environment=PORT=3001
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
StandardOutput=append:/var/log/m3u8-emitter.log
StandardError=append:/var/log/m3u8-emitter-error.log

# Seguridad adicional
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
